# etc/supervisor/conf.d/rsyncd.conf
[program:rsyncd]
command=/usr/bin/rsync --daemon --no-detach -v

stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
